# builder/alarm-builder.timer
[Unit]
Description=Build ArchLinuxARM Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target